# linux_version/systemd/wallapop-agent.timer
[Unit]
Description=Run Wallapop Agent every 2 hours

[Timer]
OnBootSec=5min
OnUnitActiveSec=120min
Persistent=true

[Install]
WantedBy=timers.target
